# rootfs/etc/services.d/tomcat/run
#!/usr/bin/with-contenv sh

echo "Starting Tomcat."
# Unable to Drop-Privs during refactor. Some image work fine, others not.
# s6-setuidgid tomcat catalina.sh run
catalina.sh run